# autostart/autostart-vpn.desktop
[Desktop Entry]
Name=autostart_vpn
Exec=/home/codayblue/.dotfiles/scripts/connect-vpn.sh <vpn uuid> > /dev/null
Terminal=false
